# files/systemd/ansible-pull-firstboot.service
[Unit]
Description=First Boot Ansible Pull Configuration
Documentation=https://docs.ansible.com/
After=network-online.target
Wants=network-online.target
Before=getty@tty1.service
ConditionPathExists=!/var/lib/ansible-pull-configured

[Service]
Type=oneshot
User=root
WorkingDirectory=/tmp
ExecStart=/usr/local/bin/bootstrap.sh
ExecStartPost=/bin/touch /var/lib/ansible-pull-configured
StandardOutput=file:/var/log/ansible-pull-firstboot.log
StandardError=file:/var/log/ansible-pull-firstboot.log
TimeoutSec=1800
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target