# files/systemd/ansible-pull.service
[Unit]
Description=Ansible Pull Service
Documentation=https://docs.ansible.com/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/tmp
ExecStart=/bin/bash -c 'BRANCH=$(/usr/local/bin/get-branch.sh); /usr/bin/ansible-pull -U https://github.com/oduvan/webquiz-ansible.git -C "$BRANCH" --clean site.yml'
StandardOutput=file:/mnt/data/ansible-pull.log
StandardError=file:/mnt/data/ansible-pull.log
TimeoutSec=1200
# Performance optimizations
Environment="ANSIBLE_PIPELINING=True"
Environment="ANSIBLE_FORKS=20"
Environment="ANSIBLE_GATHERING=smart"

[Install]
WantedBy=multi-user.target