# files/systemd/ansible-pull.service
[Unit]
Description=Ansible Pull Service
Documentation=https://docs.ansible.com/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/tmp
ExecStart=/usr/local/bin/ansible-pull-wrapper.sh
StandardOutput=file:/mnt/data/ansible-pull.log
StandardError=file:/mnt/data/ansible-pull.log
TimeoutSec=infinity
# Performance optimizations
Environment="ANSIBLE_PIPELINING=True"
Environment="ANSIBLE_FORKS=20"
Environment="ANSIBLE_GATHERING=smart"

[Install]
WantedBy=multi-user.target